# .flywheel/launchd/ai.zeststream.worker-deep-liveness-probe.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.zeststream.worker-deep-liveness-probe</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>set -uo pipefail; /Users/josh/.claude/skills/.flywheel/scripts/worker-deep-liveness-probe.sh --json &gt;&gt; /Users/josh/.local/state/flywheel/worker-deep-liveness-probe-runs.jsonl 2&gt;&gt; /Users/josh/.local/state/flywheel/worker-deep-liveness-probe.err.log; exit 0</string>
  </array>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>RunAtLoad</key>
  <false/>
  <key>KeepAlive</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/josh/.local/state/flywheel/worker-deep-liveness-probe.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/josh/.local/state/flywheel/worker-deep-liveness-probe.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/josh</string>
    <key>PATH</key>
    <string>/Users/josh/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
