# .flywheel/launchd/ai.zeststream.cross-orch-idle-watchtower.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.zeststream.cross-orch-idle-watchtower</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/josh/Developer/flywheel/.flywheel/scripts/cross-orch-idle-watchtower.sh</string>
    <string>run</string>
    <string>--mode</string>
    <string>nudge</string>
    <string>--apply</string>
    <string>--json</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/josh/Developer/flywheel</string>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/josh</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/Users/josh/.cargo/bin:/Users/josh/.local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/josh/.local/state/flywheel/cross-orch-idle-watchtower.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/josh/.local/state/flywheel/cross-orch-idle-watchtower.err.log</string>
</dict>
</plist>
