# .flywheel/launchd/com.zeststream.flywheel.w1-r1-cadence.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  W1 cadence (v5 forever-goal): fires R1 trauma-claim emitter hourly.
  Appends one row to .flywheel/state/r1-cadence-ledger.jsonl per hour.

  Load:
    cp .flywheel/launchd/com.zeststream.flywheel.w1-r1-cadence.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.zeststream.flywheel.w1-r1-cadence.plist

  Unload:
    launchctl bootout gui/$(id -u) ~/Library/LaunchAgents/com.zeststream.flywheel.w1-r1-cadence.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.zeststream.flywheel.w1-r1-cadence</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>/Users/josh/Developer/flywheel/.flywheel/scripts/w1-r1-cadence-tick.sh tick</string>
  </array>
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/josh/.local/state/flywheel/w1-r1-cadence.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/josh/.local/state/flywheel/w1-r1-cadence.stderr.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
